# server/copilot-server.service
[Unit]
Description=Copilot CLI Android WebSocket Server
Documentation=https://github.com/ssfdre38/copilot-android-client
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/copilot-android-client/server
ExecStart=/usr/bin/node /home/ubuntu/copilot-android-client/server/server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3002

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cpcli-server

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/copilot-android-client/server

[Install]
WantedBy=multi-user.target